# service.conf/docker/docker.service.d/docker-opts.conf
[Service]
ExecStart=
ExecStart=/usr/bin/docker daemon -H fd:// --storage-driver=overlay --storage-opt dm.use_deferred_removal=true --storage-opt dm.use_deferred_deletion=true
